# build/packages/eStation-Apps/build/DEBIAN/postinst
#!/bin/sh

# Fix problem for gnuplot
# Should the following line be commented? Is there a bug here? (mgm - mch)
ln -sf /lib/libreadline.so.5.2 /lib/libreadline.so.4

#
# User customization
# Configuration of the Desktop
perl -pi -e "s#/usr/share/backgrounds/warty-final-ubuntu.png#/etc/e-Station/backgrounds/e-Station.png#g" /usr/share/gconf/defaults/16_ubuntu-wallpapers
perl -pi -e "s#zoom#scaled#g" /usr/share/gconf/defaults/16_ubuntu-wallpapers

# User specification customization
user_conf=/etc/e-Station/user
user_skel=/etc/skel
${user_conf}/make_user

# Add the customization for the user
ici=$(pwd)
for user in $(ls ${user_conf}|grep "home_"| cut -d "_" -f2);do
	cd ${user_skel}
# 	su ${user} -c "find -H . -maxdepth 1 -exec cp -af ${user_skel}/{} /home/${user} \;"
	cd ${user_conf}/home_${user}
	su ${user} -c "find -H . -maxdepth 1 -exec cp -af ${user_conf}/home_${user}/{} /home/${user} \;"
	su ${user} -c "gconftool-2 --set \"/desktop/gnome/lockdown/disable_user_switching\" --type boolean true"
	su ${user} -c "gconftool-2 --set \"/apps/gnome-screensaver/user_switch_enabled\" --type boolean false"
	su ${user} -c "gconftool-2 --set \"/apps/gnome-screensaver/lock_enabled\" --type boolean true"
done
cd $ici

# USER CRON
# First creation of the crontab for each user
for user in $(ls ${user_conf}|grep "cron_"| cut -d "_" -f2);do
 	if [ ! -f "/var/spool/cron/crontabs/${user}" ]; then
		# allow user to create their cron
		 crontab -u ${user} ${user_conf}/cron_${user}
        fi
done

# Modification of the crontab entry
# Principe : Keep the tasks added by the user and add specific tasks tagged with AMESD_ID_<number>
for user in $(ls ${user_conf}/cron_* | cut -d "_" -f2) ;do
    cronfile=/var/spool/cron/crontabs/${user}
    # keep only lines not content tag with AMESD_ID_<number>
    echo "$(cat ${cronfile} | tail -n $(($(cat ${cronfile} |wc -l)-3)) | grep -v 'AMESD_ID')" > /tmp/${user}_tmp

    cat ${user_conf}/cron_${user} >> /tmp/${user}_tmp

    crontab -u ${user} /tmp/${user}_tmp
done

#
# RIGHTS
# change owner of /var/www
chown -R www-data:www-data /var/www

# NEW CONFIG FOR MUNIN (mch)
#sed -i "s@/var/cache/munin/www@/var/munin/www@g" /etc/munin/munin.conf
#mv /var/cache/munin/www/ /var/www/munin/
#chown munin.munin -R /var/www/munin
#sed -i "s@/var/cache/munin/www@/var/munin/www@g" /etc/munin/apache.conf
# END NEW CONF FOR MUNIN

chmod 777 /var/log/munin /var/www/munin /usr/share/munin/munin-graph
rm -f /var/www/index.html
ln -sf /bin/bash /bin/sh
ln -sf /usr/bin/xmlstarlet /usr/bin/xml
cp -f /etc/e-Station/freshclam.conf /etc/clamav/freshclam.conf
cp -f /etc/e-Station/sudoers /etc/sudoers
cp -f /etc/e-Station/gdmflexiserver /usr/bin/gdmflexiserver

# CLAMAV antivirus
update-rc.d -f clamav-freshclam remove

# Allow rsync between PS and EMMA station
perl -pi -e "s#RSYNC_ENABLE=false#RSYNC_ENABLE=true#" /etc/default/rsync
#/etc/init.d/rsync restart
service rsync restart

# remove mail at each reboot
#if [ ! "$(grep mail /etc/init.d/webmin 2>/dev/null)" ]; then
#	perl -pi -e "s#'start'\)#'start'\)\n	rm -f /var/spool/mail/*#g" /etc/init.d/webmin
#fi

# Copy the new php.ini (see #2294) - also comment the lines below to turn on mapscript
if [ -f /etc/e-Station/config/php.ini ]; then
	if [ -f /etc/php5/apache2/php.ini ]; then 
        mv /etc/php5/apache2/php.ini /etc/php5/apache2/php.ini.before_1.3.1
    fi
	cp -f /etc/e-Station/config/php.ini /etc/php5/apache2/php.ini
fi
# Add mapscript in PHP (see #2294)
#if [ ! "$(grep mapscript /etc/php5/apache2/php.ini 2>/dev/null)" ]; then#
#	perl -pi -e "s#^;   extension=msql.so#;   extension=msql.so\n   extension=php_mapscript.so#g" /etc/php5/apache2/php.ini
#fi
# Add phppgadmin in apache
if [ -f /etc/phppgadmin/apache.conf ]; then
	#ln -sf /etc/phppgadmin/apache.conf /etc/apache2/conf.d/phppgadmin.conf
	perl -pi -e "s|^allow from 127.0.0.0/255.0.0.0|#allow from 127.0.0.0/255.0.0.0|" /etc/phppgadmin/apache.conf
	perl -pi -e "s/# allow from all/allow from all/" /etc/phppgadmin/apache.conf
	/etc/init.d/apache2 reload #service apache2 reload #/etc/init.d/apache2 reload
fi

# MUNIN
if [ ! "$(grep "192.168.0" /etc/munin/munin.conf 2>/dev/null)" ]; then
	perl -pi -e "s#^\[localhost\.localdomain\]#[eStation-PS]\n    address 192.168.0.5\n    use_node_name yes#" /etc/munin/munin.conf
	perl -pi -e "s#^    address 127.0.0.1#[eStation-EMMA]\n    address 192.168.0.6#" /etc/munin/munin.conf
fi
if [ ! "$(grep "192\\\.168\\\.0" /etc/munin/munin-node.conf 2>/dev/null)" ]; then
	echo "allow ^192\.168\.0\.5$" >> /etc/munin/munin-node.conf
	echo "allow ^192\.168\.0\.6$" >> /etc/munin/munin-node.conf
	service munin-node restart #/etc/init.d/munin-node restart
fi
if [ -f /etc/munin/apache.conf ]; then
	perl -pi -e "s/Allow from localhost 127.0.0.0\/8 ::1/Allow from All/" /etc/munin/apache.conf
	/etc/init.d/apache2 reload #service apache2 reload #/etc/init.d/apache2 reload
fi


# Add crontab for root user for update depot
if [ ! "$(grep debdepot /var/spool/cron/crontabs/root 2>/dev/null)" ]; then
	echo "0 0 * * * /debdepot/update_depot.sh /update /debdepot > /dev/null 2>&1 #Cron for Debdepot" >> /var/spool/cron/crontabs/root
fi
if [ ! "$(grep rsync /var/spool/cron/crontabs/root 2>/dev/null)" ]; then
	echo "#*/30 * * * * /usr/bin/rsync_data data > /dev/null 2>&1 #Cron for synchronisation" >> /var/spool/cron/crontabs/root
fi
if [ ! "$(grep DumpDatabase /var/spool/cron/crontabs/root 2>/dev/null)" ]; then
	echo "0 0 * * * /usr/bin/DumpDatabase.sh  > /dev/null 2>&1 # cron for the dump" >> /var/spool/cron/crontabs/root
fi
if [ ! "$(grep clean_emma_tmp /var/spool/cron/crontabs/root 2>/dev/null)" ]; then
	echo "0 * * * * /opt/facilities/EMMA/emma/clean_emma_tmp.sh > /dev/null 2>&1 # clean tmp for EMMA" >> /var/spool/cron/crontabs/root
fi
# Mount /eumetcast 
mkdir -p -m 777 /Acquisition
if [ ! -h /eumetcast ];then
	ln -sf /Acquisition/eumetcast /eumetcast
fi
# correct bug 2261
if [ ! -h /update ];then
	ln -sf /Acquisition/update /update
fi
mkdir -p -m 777 /data_mirror
if [ ! "$(grep data_mirror /etc/auto.master 2>/dev/null)" ]; then
	echo "/data_mirror /etc/auto.estation.nfs --ghost --timeout=60"	>> /etc/auto.master
fi
if [ ! "$(grep Acquisition /etc/auto.master 2>/dev/null)" ]; then
	echo "/Acquisition /etc/auto.estation.smb  --ghost,--timeout=30" >> /etc/auto.master
fi
service autofs reload #/etc/init.d/autofs reload
if [ ! "$(grep eumetcast /etc/fstab 2>/dev/null)" ]; then
	echo "/dev/sda5	/data	ext3	defaults	0	0"			>> /etc/fstab
fi
if [ ! "$(grep "file:/debdepot" /etc/apt/sources.list 2>/dev/null)" ]; then
	echo "deb file:/debdepot lucid extra" > /etc/apt/sources.list
fi
# Samba
if [ ! "$(grep "\[data\]" /etc/samba/smb.conf 2>/dev/null)" ]; then
	echo "[data]" 			>> /etc/samba/smb.conf
	echo "	path = /data" 		>> /etc/samba/smb.conf
	echo "	writeable = yes"	>> /etc/samba/smb.conf
	echo "	browseable = yes" 	>> /etc/samba/smb.conf
	echo "	guest ok = yes" 	>> /etc/samba/smb.conf
fi

if [ ! "$(grep Acquisition /etc/hosts 2>/dev/null)" ]; then
	echo "192.168.0.2	Acquisition.eStation	Acquisition" 	>> /etc/hosts
fi
if [ ! "$(grep eStation-VIP /etc/hosts 2>/dev/null)" ]; then
	echo "192.168.0.20	eStation-VIP.eStation	eStation-VIP" 	>> /etc/hosts
fi
# Custom Firefox homepage
#CHANGE TO HAVE firefox FOLDER
ln -s /usr/lib/firefox-[0-9]* /usr/lib/firefox
#END CHANGE
# Following lines should be checked: Firefox on lucid does not have browserconfig.properties
#if [ ! "$(grep "browser.startup.homepage" /etc/xul-ext/ubufox.js 2>/dev/null)" ]; then

#echo "browser.startup.homepage=http://$(uname -n)" > /usr/lib/browserconfig.properties
#echo "browser.startup.homepage_reset=http://$(uname -n)" >> /usr/lib/browserconfig.properties

#echo "pref(\"browser.startup.homepage\",\"http://$(uname -n)\");" >> /usr/lib/firefox/defaults/pref/vendor.js
#echo "pref(\"browser.startup.homepage_reset\",\"http://$(uname -n)\");" >> /usr/lib/firefox/defaults/pref/vendor.js

#echo "pref(\"browser.startup.homepage\",\"http://$(uname -n)\");" >> /usr/lib/firefox/defaults/syspref/firefox.js
#echo "pref(\"browser.startup.homepage_reset\",\"http://$(uname -n)\");" >>  /usr/lib/firefox/defaults/syspref/firefox.js

#echo "user_pref(\"browser.startup.homepage\",\"http://$(uname -n)\");" >>  /etc/xul-ext/ubufox.js
#echo "user_pref(\"browser.startup.homepage_reset\",\"http://$(uname -n)\");" >>  /etc/xul-ext/ubufox.js
#echo "pref(\"browser.startup.homepage\",\"http://$(uname -n)\");" >>  /etc/xul-ext/ubufox.js
#echo "pref(\"browser.startup.homepage_reset\",\"http://$(uname -n)\");" >>  /etc/xul-ext/ubufox.js

echo "browser.startup.homepage=http://$(uname -n)" > /etc/xul-ext/browserconfig.properties
echo "browser.startup.homepage_reset=http://$(uname -n)" >> /etc/xul-ext/browserconfig.properties
echo "pref(\"browser.startup.homepage\",\"file:/etc/xul-ext/browserconfig.properties\");" >>  /etc/xul-ext/ubufox.js
echo "pref(\"browser.startup.homepage_reset\",\"file:/etc/xul-ext/browserconfig.properties\");" >>  /etc/xul-ext/ubufox.js


#fi

#From amesd-extern
sed -i "s|EnableProfileMigrator=.*|EnableProfileMigrator=0|" /usr/lib/firefox/application.ini
sed -i "s|EnableExtensionManager=.*|EnableExtensionManager=0|" /usr/lib/firefox/application.ini


sed -i "s#eStation\"#$(uname -n)\"#" /etc/e-Station/ubuntu-mods.js
cp -f /etc/e-Station/ubuntu-mods.js /usr/lib/firefox/extensions/ubufox@ubuntu.com/defaults/preferences/ubuntu-mods.js
# Desktop shortcut
for shortcut in upsmanager.desktop remoteaccess.desktop;do
	su adminuser -c "cp /etc/e-Station/user/home_adminuser/.local/share/applications/${shortcut} /home/adminuser/Desktop"
done

if [ "$(uname -n)" = "eStation-PS" ];then
# Modification done for PS station
	if [ ! "$(grep eStation-EMMA /etc/hosts 2>/dev/null)" ]; then
		echo "192.168.0.6	eStation-EMMA.eStation	eStation-EMMA" 	>> /etc/hosts
	fi
	update-rc.d vip defaults
	if [ ! "$(grep eStation /etc/exports 2>/dev/null)" ]; then
		echo "/data     eStation-EMMA(rw,subtree_check,async,all_squash)" >> /etc/exports
	fi
	if [ ! "$(grep eStation-EMMA /etc/auto.estation.nfs 2>/dev/null)" ]; then
		sed -i "s#eStation#eStation-EMMA#" /etc/auto.estation.nfs
	fi
	rm -f /data_processing
	ln -sf /data /data_processing
	if [ ! "$(grep emmauser /etc/gdm/custom.conf 2>/dev/null)" ]; then
		sed -i "s/\[greeter\]/\[greeter\]\nExclude=nobody,emmauser/" /etc/gdm/custom.conf
	fi
	if [ "$(grep Exclude /etc/gdm/custom.conf 2>/dev/null)" ]; then
		sed -i "s/Exclude=.*/Exclude=nobody,emmauser,ingestuser/" /etc/gdm/custom.conf
	fi
	su adminuser -c "cp /home/adminuser/.config/menus/ps_applications.menu /home/adminuser/.config/menus/applications.menu"
	for file in emma_switch-version emmauser-cron emma_recovery;do
		mv /home/adminuser/.local/share/applications/${file}.desktop /home/adminuser/.local/share/applications/${file}.desktop.recov
	done

elif [ "$(uname -n)" = "eStation-EMMA" ];then
# Modification done for EMMA station
	if [ ! "$(grep eStation-PS /etc/hosts 2>/dev/null)" ]; then
		echo "192.168.0.5	eStation-PS.eStation	eStation-PS" 	>> /etc/hosts
	fi
	if [ ! "$(grep eStation /etc/exports 2>/dev/null)" ]; then
		echo "/data     eStation-PS(rw,subtree_check,async,all_squash)" >> /etc/exports
	fi
	if [ ! "$(grep eStation-PS /etc/auto.estation.nfs 2>/dev/null)" ]; then
		sed -i "s#eStation#eStation-PS#" /etc/auto.estation.nfs
	fi
	rm -f /data_processing
	ln -sf /data_mirror/data /data_processing
	sed -i "s|^0 0 \* \* \* /usr/bin/DumpDatabase.sh|#0 0 \* \* \* /usr/bin/DumpDatabase.sh|" /var/spool/cron/crontabs/root

	if [ ! "$(grep ingestuser /etc/gdm/custom.conf 2>/dev/null)" ]; then
		sed -i "s/\[greeter\]/\[greeter\]\nExclude=nobody,deriveduser,ingestuser/" /etc/gdm/custom.conf
	fi

	su adminuser -c "cp /home/adminuser/.config/menus/emma_applications.menu /home/adminuser/.config/menus/applications.menu"
	for file in changeloglevel cleandatabase createreport exportingest importingest eumetcast-restart ingestserver-restart manualrsync changeroi createarchive eumetcast-start eumetcast-stop exportchain importchain ingestarchive ingestserver-start ingestserver-stop rsyncactivation ps_switch-version MakeDerivedActivation ps_recovery deriveduser-cron ingestuser-cron ReprocessProducts RecoverDatabase;do
		mv /home/adminuser/.local/share/applications/${file}.desktop /home/adminuser/.local/share/applications/${file}.desktop.recov
	done
	rm -f /var/spool/cron/crontabs/deriveduser /var/spool/cron/crontabs/ingestuser
	service cron restart #/etc/init.d/cron restart
fi
# Bug in the /etc/exports file replace all_squash by no_root_squash
sed -i "s|all_squash|no_root_squash|g" /etc/exports

# Restart all services modify
service autofs restart #/etc/init.d/autofs restart
service nfs-kernel-server restart #/etc/init.d/nfs-kernel-server restart
#/etc/init.d/nfs-common restart #Useless?

if [ ! "$(grep 192.168.0.1 /etc/resolv.conf 2>/dev/null)" ]; then
	echo "nameserver 192.168.0.1" >> /etc/resolv.conf
fi

# remove beep system
if [ ! "$(grep pcspkr /etc/rc.local 2>/dev/null)" ]; then
	perl -pi -e "s#^exit 0#rmmod pcspkr\nexit 0#g" /etc/rc.local
fi

echo "Modification of the Configuration"
perl -pi -e "s#data_directory = '/var/lib/postgresql/8.4/main'#data_directory = '/data/pgsql/data'#" /etc/postgresql/8.4/main/postgresql.conf
perl -pi -e "s/#listen_addresses = 'localhost'/listen_addresses = '*'/" /etc/postgresql/8.4/main/postgresql.conf
perl -pi -e "s/ssl = true/ssl = false/" /etc/postgresql/8.4/main/postgresql.conf

#Modification of the Configuration
perl -pi -e "s#host    all         all         127.0.0.1/32          md5#host    all         all         127.0.0.1/32          trust#" /etc/postgresql/8.4/main/pg_hba.conf
perl -pi -e "s#host    all         all         ::1/128               md5#host    all         all         ::1/128               trust\nhost    all         all         192.168.0.0/24               trust#" /etc/postgresql/8.4/main/pg_hba.conf 

# Update index.php (using index_1.3.0.php - see #1141)
if [ -f /var/www/index_1.3.0.php ]; then
    cp /var/www/index_1.3.0.php /var/www/index.php
fi

# For debug purposes:
# From install.log: FATAL: Could not load /lib/modules/2.6.32-38-generic/modules.dep: No such file or directory
# Because this file is under /lib/modules/2.6.32-38-generic-pae/modules.dep
# So calling again here this command maybe could left to system enough time to point to the right kernel version
# service nfs-kernel-server restart #/etc/init.d/nfs-kernel-server restart

#java java
update-alternatives --install /usr/bin/java java /usr/lib/jvm/jre1.7.0_40/bin/java 1
chmod -R a+rx /usr/lib/jvm/jre1.7.0_40/
ln -s /usr/lib/jvm/jre1.7.0_40/lib/i386/libnpjp2.so /usr/lib/firefox/plugins/


